# syssrv.service.in
[Unit]
Description=syssrv

[Service]
Type=forking
PIDFile=@app_rundir@/syssrv.pid
ExecStart=@prefix@/sbin/syssrv
ExecReload=/bin/kill -HUP $MAINPID
ExecStartPre=/usr/bin/install -d -o @app_user@ -g @app_group@ @app_rundir@ @app_logdir@

[Install]
WantedBy=multi-user.target

